# 90-keyboard_exists.rules

SUBSYSTEM=="input", RUN+="/usr/bin/udev_keyboard_exists"

